# deploy/arcade-agent.service
[Unit]
Description=Arcade CI/CD Platform Agent
Documentation=https://github.com/go-arcade/arcade
After=network.target
Wants=network.target

[Service]
Type=simple
User=arcade-agent
Group=arcade-agent
WorkingDirectory=/opt/arcade-agent
ExecStart=/opt/arcade-agent/arcade-agent -conf /opt/arcade-agent/conf.d/agent.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arcade-agent

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/arcade-agent

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 环境变量
Environment="GOMAXPROCS=0"
Environment="TZ=Asia/Shanghai"

[Install]
WantedBy=multi-user.target
